# root/etc/services.d/dhcpd/run
#!/usr/bin/with-contenv bash

/usr/sbin/dhcpd -4 -f -d --no-pid -cf /etc/dhcp/dhcpd.conf

